# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:login-server]
command=/home/RTK/src/login-server
directory=/home/RTK/src
priority=10
autostart=true
autorestart=true
startretries=10
stdout_logfile=/home/RTK/src/logs/login-stdout.log
stderr_logfile=/home/RTK/src/logs/login-stderr.log

[program:char-server]
command=/home/RTK/src/char-server
directory=/home/RTK/src
priority=20
autostart=true
autorestart=true
startretries=10
stdout_logfile=/home/RTK/src/logs/char-stdout.log
stderr_logfile=/home/RTK/src/logs/char-stderr.log

[program:map-server]
command=/home/RTK/src/map-server
directory=/home/RTK/src
priority=30
autostart=true
autorestart=true
startretries=10
stopwaitsecs=30
stdout_logfile=/home/RTK/src/logs/map-stdout.log
stderr_logfile=/home/RTK/src/logs/map-stderr.log

[program:backup]
command=/bin/bash -c "sleep 60 && while true; do /home/RTK/backup.sh; sleep 300; done"
directory=/home/RTK
priority=5
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
